# private_dot_config/systemd/user/clash.service
[Unit]
Description=clash@proxy written in go
After=network.target

[Service]
Type=simple
ExecStart=%h/bin/clash

[Install]
WantedBy=default.target
